# com.kindle.autosync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kindle.autosync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/cederic/projects/Archive/KindleClippings/kindle_auto_sync.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/cederic/projects/Archive/KindleClippings</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/cederic/projects/Archive/KindleClippings/kindle_auto_sync.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/cederic/projects/Archive/KindleClippings/kindle_auto_sync.stderr.log</string>
</dict>
</plist>